# motion-video0/root/etc/periodic/15min/clean-cameras.sh
find /var/www/localhost/htdocs/cameras -type f -ctime +1 -print | while read; do rm -f ; done
